# templates/service.erb
#!/bin/sh
#
# Sinopia deamon control script
#
# chkconfig: - 80 20
# description: Verdaccio NPM Proxy/Cache Server
# processname: Verdaccio
# pidfile: <%= @install_path %>/srv.pid
# config: <%= @install_path %>/config.yaml
Verdaccio_BIN=<%= @install_path %>/node_modules/verdaccio/bin/verdaccio
Verdaccio_USER=<%= @deamon_user %>
Verdaccio_LOG=<%= @install_path %>/deamon.log
Verdaccio_PIDFILE=<%= @install_path %>/srv.pid
EXEC_DIR=<%= @install_path %>
STOP_TIMEOUT=10
prog=Verdaccio

ECHO_CMD=/bin/echo
RM_CMD=/bin/rm
SU_CMD=/bin/su
CAT_CMD=/bin/cat
KILL_CMD=/bin/kill

start() {
		if [ -f $Verdaccio_PIDFILE ]; then
		    read ppid < $Verdaccio_PIDFILE
		    if [ `ps --pid $ppid 2> /dev/null | grep -c $ppid 2> /dev/null` -eq '1' ]; then
		      $ECHO_CMD -n "$prog is already running"
		      failure
		      $ECHO_CMD
		      return 1 
		    else
		      $RM_CMD -f $Verdaccio_PIDFILE
		    fi
		  fi
        $ECHO_CMD -n $"Starting $prog: "
        	$SU_CMD - $Verdaccio_USER -c "cd $EXEC_DIR ; nohup $Verdaccio_BIN >>$Verdaccio_LOG 2>&1 & echo \$! >$Verdaccio_PIDFILE"        	
        return $?
}

stop() {
        $ECHO_CMD -n $"Stopping $prog: "
        $KILL_CMD $($CAT_CMD ${Verdaccio_PIDFILE})
        RETVAL=$?
        $ECHO_CMD
        [ $RETVAL = 0 ] && rm -f ${Verdaccio_PIDFILE}
}

status() {
  if [ -f $Verdaccio_PIDFILE ]; then
    read ppid < $Verdaccio_PIDFILE
    if [ `ps --pid $ppid 2> /dev/null | grep -c $ppid 2> /dev/null` -eq '1' ]; then
      $ECHO_CMD "$prog is running (pid $ppid)"
      return 0
    else
      $ECHO_CMD "$prog dead but pid file exists"
      return 1
    fi
  fi
  $ECHO_CMD "$prog is not running"
  return 3
}

case "$1" in
  start)
      start
      ;;
  stop)
      stop
      ;;
  restart)
      $0 stop
      $0 start
      ;;
  status)
      status
      ;;
  *)
      ## If no parameters are given, print which are avaiable.
      $ECHO_CMD "Usage: $0 {start|stop|status|restart}"
      exit 1
      ;;
esac
